# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:collectord]
command=collectord -vv -c /etc/collectord/collectord.conf
startsecs = 0
exitcodes = 0
